# speedtest.crontab
*   *   *   *   *   echo "Running speed test." && (speedtest-cli --json >> /tests/`date +'%s'`.json) && echo "Speed test completed."
